# resq/deployment/supervisor/resq-worker.conf
; ResQ Laravel Queue Worker Supervisor Configuration
; Place this file in /etc/supervisor/conf.d/ on Ubuntu/Debian systems
; or /etc/supervisord.d/ on RHEL/CentOS systems

[program:resq-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/resq/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --queue=default,notifications,ai-chat
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/resq/storage/logs/worker.log
stopwaitsecs=3600

; Environment variables
environment=Laravel="/var/www/resq"
